# scripts/gentoo/init.d/tigase
#!/sbin/runscript

depend() {
	need net
	use dns logger postgresql mysql
}

start()	{
	ebegin "Starting Tigase"
	su - ${TIGASE_USER} -c "${TIGASE_ENV} ${TIGASE_HOME}/bin/tigase.sh start ${TIGASE_HOME}/${TIGASE_CONF}"
	eend $?
}

stop ()	{
	ebegin "Stopping Tigase"
	su - ${TIGASE_USER} -c "${TIGASE_ENV} ${TIGASE_HOME}/bin/tigase.sh stop ${TIGASE_HOME}/${TIGASE_CONF}"
	sleep 2
	eend $?
}

svc_restart () {
	ebegin "Restarting Tigase"
	su - ${TIGASE_USER} -c "${TIGASE_ENV} ${TIGASE_HOME}/bin/tigase.sh stop ${TIGASE_HOME}/${TIGASE_CONF}"
	sleep 10
	su - ${TIGASE_USER} -c "${TIGASE_ENV} ${TIGASE_HOME}/bin/tigase.sh start ${TIGASE_HOME}/${TIGASE_CONF}"
	eend $?
}
